# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/1910_train/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.test >/var/log/1910_train/crontab.log 2>&1
